# deploy/systemd/sandalpunk.service
[Unit]
Description=Sandalpunk Telegram Mini App
After=network.target

[Service]
Type=simple
User=sandalpunk
WorkingDirectory=/opt/sandalpunk
EnvironmentFile=/opt/sandalpunk/.env
ExecStart=/usr/bin/java -jar /opt/sandalpunk/sandalpunk.jar
SuccessExitStatus=143
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

